# systemd/system/paccache.service.d/invocation.conf
[Service]
ExecStart=
ExecStart=/usr/bin/paccache -rk2
ExecStart=/usr/bin/paccache -ruk0
